# supervisord.conf
[supervisord]
nodaemon=true

[program:cake_queue_worker]
command=bin/cake queue runworker
directory=/var/www/html/app
autostart=true
autorestart=true
startsecs=10
startretries=3
